# guapow/dist/daemon/systemd/user/guapow-opt.service
[Unit]
Description=guapow-opt

[Service]
ExecStart=/usr/bin/guapow-opt
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
